# birb.timer
[Unit]
Description = run birb.py minutely

[Timer]
OnCalendar = minutely

[Install]
Wantedby = timers.target